# friendheat/files/friendheat-tornado.conf
description "daemon for friendheat"

stop on runlevel[06]
start on started nginx

console output
exec /home/stillinbeta/sites/friendheat/bin/python /home/stillinbeta/sites/friendheat/application.py
