# statusbot/templates/upstart.conf.erb
description "Statusbot"
author  "Ellie Frost <web@stillinbeta.com>"

start on runlevel [234]
stop on runlevel [0156]

setuid stillinbeta
chdir <%= scope['statusbot::site_dir'] %>
exec bin/python /usr/bin/twistd --syslog -ny statusbot.tac
respawn
