# statusbot/templates/upstart.conf.erb
description "Statusbot"
author  "Ellie Frost <web@stillinbeta.com>"

start on runlevel [234]
stop on runlevel [0156]

chdir <%= site_dir %>
exec bin/twistd --syslog --uid=stillinbeta -ny statusbot.tac
respawn
