# gunicorn/templates/django_upstart.conf.erb
description "gunicorn for <%= name %>"

stop on runlevel[06]
start on started nginx

console output
exec <%= directory %>/bin/gunicorn \
     --bind=unix:/tmp/gunicorn-<%= site_name %>.sock \
     --error-log=/tmp/<%= site_name %>-error.log \
     --pythonpath=<%= directory %> \
     django-wsgi:application

